# resiliosync/rootfs/etc/s6-overlay/s6-rc.d/resiliosync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# barisozdag's Personal Home Assistant Add-ons: Resilio Sync
# Runs the Resilio Sync
# ==============================================================================
bashio::log.info 'Starting the Resilio Sync...'

# Run the Resilio Sync
exec rslsync --nodaemon --config /config/resiliosync/sync.conf
